# systemd/system/snap.nextcloud.apache.service
[Unit]
# Auto-generated, DO NOT EDIT
Description=Service for snap application nextcloud.apache
Requires=snap-nextcloud-6916.mount
Wants=network-online.target
After=snap-nextcloud-6916.mount network-online.target
X-Snappy=yes

[Service]
ExecStart=/usr/bin/snap run nextcloud.apache
SyslogIdentifier=nextcloud.apache
Restart=always
WorkingDirectory=/var/snap/nextcloud/6916
ExecStop=/usr/bin/snap run --command=stop nextcloud.apache
TimeoutStopSec=30
Type=simple

[Install]
WantedBy=multi-user.target
